# files/etc/systemd/system/vdisk-relay-web-status-cache.service
[Unit]
Description=vdisk-relay web status cache

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vdisk-relay-web-status-cache
Nice=10
