# files/etc/systemd/system/vdisk-relay-status-cache.service
[Unit]
Description=vdisk-relay cached login status

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vdisk-relay-status-cache
Nice=10
