# files/etc/systemd/system/vdisk-relay-debug-cache.service
[Unit]
Description=vdisk-relay debug status cache

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vdisk-relay-debug-cache
Nice=10
