# files/etc/systemd/system/vdisk-relay-health-live.service
[Unit]
Description=vdisk-relay live USB/image stale watchdog
ConditionPathExists=/usr/local/sbin/vdisk-relay
ConditionPathExists=/usbdisk.img
ConditionPathExists=/etc/vdisk-relay.live-enabled
ConditionPathExists=/etc/vdisk-relay.allow-health-reboot

[Service]
Type=oneshot
ExecStart=/bin/bash -lc '/usr/local/sbin/vdisk-relay image-watchdog --send --reboot-stale || true'
Nice=10
