# files/etc/systemd/system/vdisk-relay-fast-sync.service
[Unit]
Description=vdisk-relay fast live sync after pending
After=local-fs.target
ConditionPathExists=/usr/local/sbin/vdisk-relay
ConditionPathExists=/usbdisk.img
ConditionPathExists=/etc/vdisk-relay.live-enabled
ConditionPathExists=/etc/vdisk-relay.allow-sync

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vdisk-relay fast-sync
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
