# files/etc/systemd/system/vdisk-relay-sync-live.service
[Unit]
Description=vdisk-relay live sync once
After=local-fs.target
ConditionPathExists=/usr/local/sbin/vdisk-relay
ConditionPathExists=/usbdisk.img
ConditionPathExists=/etc/vdisk-relay.live-enabled
ConditionPathExists=/etc/vdisk-relay.allow-sync

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/vdisk-relay live-sync --force
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
