# files/etc/update-motd.d/99-vdisk-relay
#!/bin/sh

[ -r /run/vdisk-relay/status-login.txt ] || exit 0

echo
cat /run/vdisk-relay/status-login.txt
echo
